# latest/overlay/etc/systemd/system/ipfs.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/ipfs daemon
User=root

[Install]
WantedBy=multiuser.target
